# cookbooks/supervisor/templates/default/program.conf.erb
[program:<%= @prog.pattern %>]
command=<%= @prog.pattern %>