# cookbooks/rails/templates/default/app_supervisor.conf.erb
[program:rails-<%= @name %>]
command = /bin/bash --login -c 'unicornherder <%= node['rails']['applications_root'] %>/<%= @name %>/current/bin/unicorn -- <%= node['rails']['applications_root'] %>/<%= @name %>/current/config.ru -Dc <%= node['rails']['applications_root'] %>/<%= @name %>/shared/config/unicorn.rb -E <%= @rails_env %>'
user = deploy